# pkg/rpm/filelist
%attr(- root root) /etc/gbootroot/gbootrootrc
%attr(- root root) /usr/bin/gbootroot
%attr(- root root) /usr/bin/linuxbr
%attr(- root root) /usr/bin/uml_mconsole
%attr(- root root) /usr/bin/uml_moo
%attr(- root root) /usr/bin/uml_net
%attr(- root root) /usr/bin/tunctl
%attr(- root root) /usr/bin/make_debian
%attr(- root root) /usr/bin/uml_switch
%attr(- root root) /usr/lib/uml/port-helper
%attr(- root root) /usr/lib/bootroot/yard/Replacements/lib/modules/modules-2.2.tar
%attr(- root root) /usr/lib/bootroot/yard/Replacements/lib/modules/modules-2.4.tar
%attr(- root root) /usr/lib/bootroot/yard/Replacements/lib/modules/config
%attr(- root root) /usr/lib/bootroot/expect_uml
%attr(- root root) /usr/lib/bootroot/yard_chrooted_tests
%attr(- root root) /usr/lib/bootroot/genext2fs
%attr(- root root) /usr/lib/bootroot/root_filesystem/Initrd.gz
%attr(- root root) /usr/lib/bootroot/root_filesystem/root_fs_helper
%attr(- root root) /usr/X11R6/include/X11/pixmaps/gbootroot.xpm
%attr(- root root) /usr/share/doc/gbootroot/html/index.html
%attr(- root root) /usr/share/doc/gbootroot/html/images/gbootroot.jpg
%attr(- root root) /usr/share/doc/gbootroot/html/images/Editorspick.gif
%attr(- root root) /usr/share/doc/gbootroot/html/images/peng-movie.4.gif
%attr(- root root) /usr/share/doc/gbootroot/html/images/paths.jpg
%attr(- root root) /usr/share/doc/gbootroot/html/images/create.jpg
%attr(- root root) /usr/share/doc/gbootroot/html/images/template_search.jpg
%attr(- root root) /usr/share/doc/gbootroot/html/images/file.jpg
%attr(- root root) /usr/share/doc/gbootroot/html/images/screenshot.jpg
%attr(- root root) /usr/share/doc/gbootroot/html/images/verbosity_box.jpg
%attr(- root root) /usr/share/doc/gbootroot/html/images/ABS.jpg
%attr(- root root) /usr/share/doc/gbootroot/html/images/settings.jpg
%attr(- root root) /usr/share/doc/gbootroot/html/images/stripping.jpg
%attr(- root root) /usr/share/doc/gbootroot/html/images/screenshot2.jpg
%attr(- root root) /usr/share/doc/gbootroot/html/images/ARS.jpg
%attr(- root root) /usr/share/doc/gbootroot/html/images/replacements.jpg
%attr(- root root) /usr/share/doc/gbootroot/html/images/rateit80x18.gif
%attr(- root root) /usr/share/doc/gbootroot/html/images/gBSicon.jpg
%attr(- root root) /usr/share/doc/gbootroot/html/images/filesystem.jpg
%attr(- root root) /usr/share/doc/gbootroot/html/images/tests.jpg
%attr(- root root) /usr/share/doc/gbootroot/html/images/uml_box.jpg
%attr(- root root) /usr/share/doc/gbootroot/html/images/gBS.jpg
%attr(- root root) /usr/share/doc/gbootroot/html/images/yard_box.jpg
%attr(- root root) /usr/share/doc/gbootroot/html/bootroot.html
%attr(- root root) /usr/share/doc/gbootroot/Changes
%attr(- root root) /usr/share/gbootroot/yard/Replacements/CVS/Root
%attr(- root root) /usr/share/gbootroot/yard/Replacements/CVS/Repository
%attr(- root root) /usr/share/gbootroot/yard/Replacements/CVS/Entries
%attr(- root root) /usr/share/gbootroot/yard/Replacements/dev/CVS/Root
%attr(- root root) /usr/share/gbootroot/yard/Replacements/dev/CVS/Repository
%attr(- root root) /usr/share/gbootroot/yard/Replacements/dev/CVS/Entries
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/rc
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/CVS/Root
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/CVS/Repository
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/CVS/Entries
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/apt/CVS/Root
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/apt/CVS/Repository
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/apt/CVS/Entries
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/apt/sources.list
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/motd
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/hosts
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/pam.d/CVS/Root
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/pam.d/CVS/Repository
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/pam.d/CVS/Entries
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/pam.d/other
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/inittab.example-helper
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/inittab.mingetty
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/group.example
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/pam.conf
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/linuxrc
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/group-debian
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/ioctl.save
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/termcap
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/inittab.example-deb-nodevfs
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/group.debian
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/nsswitch.conf.example
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/passwd-debian
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/fstab.debian
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/fstab.example
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/securetty.debian
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/network/CVS/Root
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/network/CVS/Repository
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/network/CVS/Entries
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/network/interfaces
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/init.d/CVS/Root
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/init.d/CVS/Repository
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/init.d/CVS/Entries
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/init.d/rcS.helper
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/init.d/halt.example
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/init.d/reboot.example
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/init.d/rc.example
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/init.d/rcS.example
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/hostname
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/inittab.example.mingetty-rpm
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/passwd.debian
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/ttytype
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/gettydefs
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/passwd
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/inittab.example-deb
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/networks
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/ld.so.conf
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/inittab.example.mingetty-rpm-nodevfs
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/securetty.example
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/inittab.example.agetty-slack
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/passwd.example
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/inittab
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/inittab.agetty
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/inittab.example.agetty-slack-nodevfs
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/inittab.debian
%attr(- root root) /usr/share/gbootroot/yard/Replacements/etc/shadow.debian
%attr(- root root) /usr/share/gbootroot/yard/Replacements/usr/CVS/Root
%attr(- root root) /usr/share/gbootroot/yard/Replacements/usr/CVS/Repository
%attr(- root root) /usr/share/gbootroot/yard/Replacements/usr/CVS/Entries
%attr(- root root) /usr/share/gbootroot/yard/Replacements/usr/src/CVS/Root
%attr(- root root) /usr/share/gbootroot/yard/Replacements/usr/src/CVS/Repository
%attr(- root root) /usr/share/gbootroot/yard/Replacements/usr/src/CVS/Entries
%attr(- root root) /usr/share/gbootroot/yard/Replacements/usr/src/linux-um/CVS/Root
%attr(- root root) /usr/share/gbootroot/yard/Replacements/usr/src/linux-um/CVS/Repository
%attr(- root root) /usr/share/gbootroot/yard/Replacements/usr/src/linux-um/CVS/Entries
%attr(- root root) /usr/share/gbootroot/yard/Replacements/usr/src/linux-um/drivers/CVS/Root
%attr(- root root) /usr/share/gbootroot/yard/Replacements/usr/src/linux-um/drivers/CVS/Repository
%attr(- root root) /usr/share/gbootroot/yard/Replacements/usr/src/linux-um/drivers/CVS/Entries
%attr(- root root) /usr/share/gbootroot/yard/Replacements/usr/src/linux-um/drivers/mtd/CVS/Root
%attr(- root root) /usr/share/gbootroot/yard/Replacements/usr/src/linux-um/drivers/mtd/CVS/Repository
%attr(- root root) /usr/share/gbootroot/yard/Replacements/usr/src/linux-um/drivers/mtd/CVS/Entries
%attr(- root root) /usr/share/gbootroot/yard/Replacements/usr/src/linux-um/drivers/mtd/devices/CVS/Root
%attr(- root root) /usr/share/gbootroot/yard/Replacements/usr/src/linux-um/drivers/mtd/devices/CVS/Repository
%attr(- root root) /usr/share/gbootroot/yard/Replacements/usr/src/linux-um/drivers/mtd/devices/CVS/Entries
%attr(- root root) /usr/share/gbootroot/yard/Replacements/usr/src/linux-um/drivers/mtd/devices/mtdram.o
%attr(- root root) /usr/share/gbootroot/yard/Replacements/usr/src/linux-um/drivers/mtd/devices/blkmtd.o
%attr(- root root) /usr/share/gbootroot/yard/Replacements/home/CVS/Root
%attr(- root root) /usr/share/gbootroot/yard/Replacements/home/CVS/Repository
%attr(- root root) /usr/share/gbootroot/yard/Replacements/home/CVS/Entries
%attr(- root root) /usr/share/gbootroot/yard/Replacements/home/user/CVS/Root
%attr(- root root) /usr/share/gbootroot/yard/Replacements/home/user/CVS/Repository
%attr(- root root) /usr/share/gbootroot/yard/Replacements/home/user/CVS/Entries
%attr(- root root) /usr/share/gbootroot/yard/Replacements/home/user/README
%attr(- root root) /usr/share/gbootroot/yard/Replacements/home/user/.bashrc.debian
%attr(- root root) /usr/share/gbootroot/yard/Replacements/home/user/.bash_profile.debian
%attr(- root root) /usr/share/gbootroot/yard/Replacements/mnt1/CVS/Root
%attr(- root root) /usr/share/gbootroot/yard/Replacements/mnt1/CVS/Repository
%attr(- root root) /usr/share/gbootroot/yard/Replacements/mnt1/CVS/Entries
%attr(- root root) /usr/share/gbootroot/yard/Replacements/mnt2/CVS/Root
%attr(- root root) /usr/share/gbootroot/yard/Replacements/mnt2/CVS/Repository
%attr(- root root) /usr/share/gbootroot/yard/Replacements/mnt2/CVS/Entries
%attr(- root root) /usr/share/gbootroot/yard/Replacements/root/CVS/Root
%attr(- root root) /usr/share/gbootroot/yard/Replacements/root/CVS/Repository
%attr(- root root) /usr/share/gbootroot/yard/Replacements/root/CVS/Entries
%attr(- root root) /usr/share/gbootroot/yard/Replacements/root/.profile.debian
%attr(- root root) /usr/share/gbootroot/yard/Replacements/root/.profile
%attr(- root root) /usr/share/gbootroot/yard/Replacements/root/.bashrc.debian
%attr(- root root) /usr/share/gbootroot/yard/Replacements/root/umlnet
%attr(- root root) /usr/share/gbootroot/yard/templates/Example-Debian.yard
%attr(- root root) /usr/share/gbootroot/yard/templates/Example-Mini.yard
%attr(- root root) /usr/share/gbootroot/yard/templates/Example-EmptyFS.yard
%attr(- root root) /usr/share/gbootroot/yard/templates/Initrd.yard
%attr(- root root) /usr/share/gbootroot/yard/templates/Helper.yard
%attr(- root root) /usr/share/gbootroot/yard/templates/Example.yard
%attr(- root root) /usr/share/gbootroot/genext2fs/dev.txt
%attr(- root root) /usr/share/gbootroot/genext2fs/Makefile
%attr(- root root) /usr/share/gbootroot/genext2fs/device_table.txt
%attr(- root root) /usr/share/gbootroot/genext2fs/genext2fs.c
%attr(- root root) /usr/share/perl5/BootRoot/lsMode.pm
%attr(- root root) /usr/share/perl5/BootRoot/Options.pm
%attr(- root root) /usr/share/perl5/BootRoot/UML.pm
%attr(- root root) /usr/share/perl5/BootRoot/BootRoot.pm
%attr(- root root) /usr/share/perl5/BootRoot/Yard.pm
%attr(- root root) /usr/share/perl5/BootRoot/YardBox.pm
%attr(- root root) /usr/share/perl5/BootRoot/Error.pm
